# debian/renderd.service
[Unit]
Description=Renderd

[Service]
ExecStart=/etc/init.d/renderd start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Type=notify

[Install]
